# local-portal.service
[Unit]
Description=local-portal - dev port guard and registry
After=network.target

[Service]
Type=simple
# Update this path if you install local-portal somewhere else.
WorkingDirectory=%h/project/Shili/workspaces/dev-lite/local-portal
ExecStart=/usr/bin/env node src/main.js
Restart=on-failure
RestartSec=5
# exit code 2 means another local-portal instance already owns this apiPort - don't fight it.
RestartPreventExitStatus=2
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
